# ntfy/systemd/crowdsec-watch.service
[Unit]
Description=CrowdSec health watch -> ntfy
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/crowdsec-mikrotik-bouncer/ntfy/crowdsec-watch.sh
